# scripts/deploy/claude-bot.service.template
[Unit]
Description=Claude Code Multi-Platform Bot (Feishu/DingTalk/WeChat)
After=network.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
# 用户和目录由安装脚本动态配置
# User=<BOT_USER>
# WorkingDirectory=<BOT_DIR>
ExecStart=<BUN_PATH> run bot
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=30

# 环境变量
Environment="NODE_ENV=production"
Environment="BUN_VERSION=1.3.9"
Environment="WS_PORT=8765"
Environment="WEBHOOK_PORT=3000"

# 日志配置 (由安装脚本动态配置)
# StandardOutput=append:<LOG_DIR>/bot.log
# StandardError=append:<LOG_DIR>/bot-error.log

# 资源限制
LimitNOFILE=65535
MemoryMax=2G

# 安全加固
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target